# scripts/update/res/v0.0.0/hyparvisor.service
[Unit]
Description=Hyparvisor web server
After=network.target

[Service]
ExecStart=/bin/bash /opt/hyparvisor/hyparvisor.sh
Type=exec
Restart=always
StandardOutput=tty
StandardError=tty

[Install]
WantedBy=getty.target
